# deploy/memoryx-hermes.service
[Unit]
Description=Mnemosyne-X Memory Service
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/memoryx
Environment=MEMORYX_ENV_FILE=%h/memoryx/.env
ExecStart=%h/memoryx/.venv/bin/python -m memoryx
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
